# wg/qbittorrent_namespace/vpn-namespace.service
[Unit]
Description=VPN Network Namespace for qBittorrent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/vpn-namespace-setup.sh
ExecStop=/usr/local/bin/vpn-namespace-cleanup.sh

[Install]
WantedBy=multi-user.target
